# bbs.conf
[program:BBS]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/var/www/BBS
autostart=true
autorestart=true
